# teddycloud/rootfs/etc/s6-overlay/s6-rc.d/teddycloud/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: TeddyCloud
# Runs TeddyCloud
# ==============================================================================
bashio::log.info "Starting TeddyCloud..."

# Erstelle Verzeichnisse in addon_config
mkdir -p /addon_config/teddycloud/certs/server
mkdir -p /addon_config/teddycloud/certs/client
mkdir -p /addon_config/teddycloud/config
mkdir -p /addon_config/teddycloud/data/content
mkdir -p /addon_config/teddycloud/data/library
mkdir -p /addon_config/teddycloud/data/custom_img
mkdir -p /addon_config/teddycloud/data/firmware
mkdir -p /addon_config/teddycloud/data/cache

# Migration: Falls alte Daten in /config/teddycloud vorhanden
if [[ -d /config/teddycloud/certs ]] && [[ ! -L /config/teddycloud/certs ]]; then
    bashio::log.info "Migrating old certs..."
    cp -rn /config/teddycloud/certs/. /addon_config/teddycloud/certs/ 2>/dev/null || true
fi
if [[ -d /config/teddycloud/config ]] && [[ ! -L /config/teddycloud/config ]]; then
    bashio::log.info "Migrating old config..."
    cp -rn /config/teddycloud/config/. /addon_config/teddycloud/config/ 2>/dev/null || true
fi

# Symlinks von TeddyCloud-Pfaden auf addon_config
rm -rf /teddycloud/certs && ln -sf /addon_config/teddycloud/certs /teddycloud/certs
rm -rf /teddycloud/config && ln -sf /addon_config/teddycloud/config /teddycloud/config
rm -rf /teddycloud/data/content && ln -sf /addon_config/teddycloud/data/content /teddycloud/data/content
rm -rf /teddycloud/data/library && ln -sf /addon_config/teddycloud/data/library /teddycloud/data/library
rm -rf /teddycloud/data/firmware && ln -sf /addon_config/teddycloud/data/firmware /teddycloud/data/firmware
rm -rf /teddycloud/data/cache && ln -sf /addon_config/teddycloud/data/cache /teddycloud/data/cache
rm -rf /teddycloud/data/www/custom_img && ln -sf /addon_config/teddycloud/data/custom_img /teddycloud/data/www/custom_img

bashio::log.info "All paths configured, starting TeddyCloud..."
cd /teddycloud || exit
exec /usr/local/bin/teddycloud
