# teddycloud/rootfs/etc/s6-overlay/s6-rc.d/teddycloud/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: TeddyCloud
# Runs TeddyCloud
# Alle Daten werden unter /config/teddycloud/ gespeichert (addon_config:rw)
# ==============================================================================

bashio::log.info "Starting TeddyCloud..."

# Erstelle alle Verzeichnisse unter /config/teddycloud/ (addon_config)
mkdir -p /config/teddycloud/certs/server
mkdir -p /config/teddycloud/certs/client
mkdir -p /config/teddycloud/config
mkdir -p /config/teddycloud/data/content
mkdir -p /config/teddycloud/data/library
mkdir -p /config/teddycloud/data/custom_img
mkdir -p /config/teddycloud/data/firmware
mkdir -p /config/teddycloud/data/cache
mkdir -p /config/teddycloud/data/plugins

# Symlinks: TeddyCloud-Pfade → /config/teddycloud/
# certs und config direkt
rm -rf /teddycloud/certs && ln -sf /config/teddycloud/certs /teddycloud/certs
rm -rf /teddycloud/config && ln -sf /config/teddycloud/config /teddycloud/config

# data/* Unterverzeichnisse
rm -rf /teddycloud/data/content && ln -sf /config/teddycloud/data/content /teddycloud/data/content
rm -rf /teddycloud/data/library && ln -sf /config/teddycloud/data/library /teddycloud/data/library
rm -rf /teddycloud/data/firmware && ln -sf /config/teddycloud/data/firmware /teddycloud/data/firmware
rm -rf /teddycloud/data/cache && ln -sf /config/teddycloud/data/cache /teddycloud/data/cache
rm -rf /teddycloud/data/www/custom_img && ln -sf /config/teddycloud/data/custom_img /teddycloud/data/www/custom_img
rm -rf /teddycloud/data/www/plugins && ln -sf /config/teddycloud/data/plugins /teddycloud/data/www/plugins

bashio::log.info "All paths linked to /config/teddycloud/"
cd /teddycloud || exit
exec /usr/local/bin/teddycloud
